# supervisord.conf
[supervisord]
nodaemon=true

[program:tracelyzer]
command=/etc/init.d/tracelyzer start -D 

[program:node]
command=/usr/bin/node /var/node/server.js -DFOREGROUND
